# scripts/compound/launchagents/com.compound.caffeinate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.compound.caffeinate</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Keeps Mac awake for 9 hours (5 PM to 2 AM = 9 hours = 32400 seconds) -->
        <string>/usr/bin/caffeinate -i -t 32400</string>
    </array>
    
    <!-- Start at 5:00 PM daily (before the loops run) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/compound-caffeinate.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/compound-caffeinate.err.log</string>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
